# rfn-mtk-linkit-webui/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# Enable RTC clock
echo ds1307 0x68 > /sys/bus/i2c/devices/i2c-0/new_device
sleep 1
hwclock -s

# Start modem
modem_ctrl on

exit 0
