# rfn-tag-receiver/files/etc/init.d/rftagreceiverd
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 OpenWrt.org

START=98
STOP=9

SERVICE_DAEMONIZE=1
SERVICE_STOP_TIME=9

boot()
{
    start
}

start()
{
    echo "Starting RF-Networks Tag Receiver application"
    service_start /usr/bin/TagReceiver
}

stop() 
{
    echo "Stopping RF-Networks Tag Receiver application"
    service_stop /usr/bin/TagReceiver
}


