# ecoppia-station/files/etc/uci-defaults/99_add_ecoppia_to_crontab
#!/bin/sh

echo "Adding to crontab"
echo '*/1 * * * * /bin/sh /usr/bin/rtr_clnt_wd 2>&1 | /usr/bin/logger' >> /etc/crontabs/root
/etc/init.d/cron reload
