# ecoppia-station/files/etc/uci-defaults/99_set_root_password
#!/bin/sh

echo 'Setting up root user default Ecoppia password'
NEWPASSWD=coMStation100!
passwd <<EOF
$NEWPASSWD
$NEWPASSWD
EOF

return 0
